# etc/init.d/nexus
#!/bin/bash
# chmod 775
# chkconfig: 35 95 1
# description: script to start/stop nexus

export RUN_AS_USER=root
/usr/local/nexus/nexus-2.7.0-05/bin/nexus $1